# app/system/system-maintenance/setup/all/host.feature/setup/Linux/files/_ROOT_/etc/systemd/system/update-system.timer
[Unit]
Description=Update system

[Timer]
OnCalendar=0/2:00:00
Persistent=true

[Install]
WantedBy=timers.target
