# app/system/firewall/setup/FreeBSD/files/_ROOT_/usr/local/etc/periodic/daily/firewall-update-country-blocks
#!/bin/sh

firewall-update-country-blocks
