# app/system/remote/setup/Linux/files/_ROOT_/etc/systemd/system/remote-publish-public-IP.service
[Unit]
Description=Publish public IP

[Service]
User=remote
Group=remote
Type=oneshot
ExecStart=/usr/local/walterjwhite/remote/bin/publish-public-IP
