# app/system/remote/setup/Linux/files/_ROOT_/etc/systemd/system/remote-publish-public-IP.timer
[Unit]
Description=publish public IP

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
