# app/system/remote/setup/Linux/files/_ROOT_/etc/systemd/system/remote-publish-public-IP.service
[Unit]
Description=Publish public IP

[Service]
Type=oneshot
ExecStart=/usr/local/bin/remote-publish-public-IP


