# app/system/system-maintenance/setup/all/zfs.feature/setup/Linux/files/_ROOT_/etc/systemd/system/zfs-pool-maintenance-monthly.service
[Unit]
Description=ZFS pool maintenance

[Service]
Type=oneshot
ExecStart=zfs-pool-maintenance
