# app/system/system-maintenance/setup/all/zfs.feature/setup/Linux/files/_ROOT_/etc/systemd/system/zfs-zap-maintenance-weekly.service
[Unit]
Description=ZFS maintenance

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zfs-zap-maintenance weekly
