# app/system/system-maintenance/setup/all/zfs.feature/setup/Linux/files/_ROOT_/etc/systemd/system/zfs-zap-maintenance-monthly.service
[Unit]
Description=ZFS maintenance

[Service]
Type=oneshot
ExecStart=zfs-zap-maintenance monthly
