# app/system/remote/setup/FreeBSD/files/_ROOT_/usr/local/etc/periodic/daily/remote-publish-public-IP
#!/bin/sh

[ -n "$_OPTN_REMOTE_PUBLISH_PUBLIC_IP_NIGHTLY_DISABLE" ] && {
  log_warn "publish public IP is disabled"
  exit 0
}

remote-publish-public-IP
