# app/firewall/setup/all/crontab
* * * * * firewall-update-tables

0 */4 * * * firewall-update-ip-blocks
