# app/remote/files/Linux/_ROOT_/etc/cron.daily/remote-publish-public-IP
#!/bin/sh

[ -n "$_OPTN_REMOTE_PUBLISH_PUBLIC_IP_NIGHTLY_DISABLE" ] && {
  _WARN "Publish public IP is disabled"
  exit 0
}

remote-publish-public-IP
