# app/system/firewall/setup/Linux/files/_ROOT_/etc/systemd/system/firewall-update-country-blocks.service
[Unit]
Description=Update country blocks

[Service]
Type=oneshot
ExecStart=/usr/local/bin/firewall-update-country-blocks
