# app/system/system-maintenance/setup/all/host.feature/setup/FreeBSD/files/_ROOT_/usr/local/etc/periodic/daily/00.app-update
#!/bin/sh

log_logfile=/var/log/walterjwhite/app-update.log app-update
