# app/system/status/setup/Linux/files/_ROOT_/etc/systemd/system/status.service
[Unit]
Description=check status of system

[Service]
Type=oneshot
ExecStart=/usr/local/bin/status
