# app/system/display/setup/all/screen-brightness.feature/files/etc/local.d/03.screen-brightness-permissions.start
#!/bin/sh

printf '%s\n' $0

find /sys/devices/ -type f -path '*/*_backlight/brightness' -exec chmod 664 {} +
find /sys/devices/ -type f -path '*/*_backlight/brightness' -exec chgrp video {} +
