# app/system/domain-blocks/setup/all/dnsmasq.feature/setup/FreeBSD/files/_ROOT_/usr/local/etc/periodic/daily/10.1.dnsmasq.blocks
dnsmasq-update-domain-blocks
service dnsmasq restart
